# deploy/systemd/sentinel-dream.service
[Unit]
Description=Sentinel nightly consolidation (dreaming)
After=sentinel-cloud.service

[Service]
Type=oneshot
# Trigger one consolidation pass on the local cloud brain.
ExecStart=/usr/bin/curl -fsS -X POST http://127.0.0.1:8000/consolidate
