# deploy/systemd/sentinel-cloud.service
[Unit]
Description=Sentinel cloud brain (Qwen-powered guardian)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/sentinel/src
# LIVE: SENTINEL_SIM=0 uses real qwen3.7-max and requires QWEN_API_KEY in /root/sentinel/.env.
# Set to 1 for safe offline SIM mode (deterministic, no API key needed).
Environment=SENTINEL_SIM=0
ExecStart=/root/sentinel/.venv/bin/python -m uvicorn server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3
# self-healing: if it dies >5x in 60s, still keep trying (don't give up)
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
